# udev/50-backup.rules
# udev by junior
KERNEL=="sd?1", ACTION=="add", SUBSYSTEMS=="usb", ATTRS{serial}=="11111", RUN+="/bin/su root -c '/scripts/run.sh %k'"
